# templates/heartbeat.service
[Unit]
Description=Heartbeat to Riemann

[Service]
Type=simple
ExecStart=/usr/bin/python /root/riemannHeartbeat.py
